# deploy/pocketbase.service
# systemd unit for PocketBase on the EC2 instance.
#
# The data dir (/opt/attendance/pb_data) holds the SQLite DB + uploaded selfies
# and MUST live on the persistent EBS volume so the nightly stop/start loses no
# data. The pb_hooks and pb_migrations dirs ship alongside the binary.
#
# Install:
#   sudo cp pocketbase.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pocketbase

[Unit]
Description=PocketBase (Attendance backend)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pocketbase
Group=pocketbase
WorkingDirectory=/opt/attendance
# ADMS_SHARED_SECRET authorizes biometric device pushes (see adms.pb.js).
Environment=ADMS_SHARED_SECRET=change-me
ExecStart=/opt/attendance/pocketbase serve --http=127.0.0.1:8090 --dir=/opt/attendance/pb_data
Restart=always
RestartSec=5
LimitNOFILE=4096
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
